# crontab
MAILTO="michael.nolan@rescue.org"
SHELL=/bin/bash

1 1 * * * cd /app/ && source .env && python3 -m scripts.populate_mongo >> /dev/stdout 2>&1
