# crontab
MAILTO=""

* * * * * echo 'AAAAAAAAAAAAAAAAAAAAAAAAAAAAAAAAAAAAaaaaaaa' 2>&1 || true
