# crontab
MAILTO=""

0 12 * * * cd /app/ && python3 -m scripts.populate_mongo